# runit/php/run
#!/bin/sh
mkdir /run/php/
exec /usr/sbin/php-fpm7.4 --nodaemonize -R