# runit/mariadb/run
#!/bin/sh

chown -R mysql:mysql /var/lib/mysql

exec 2>&1
exec /usr/bin/mysqld_safe